# src/supervisor.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx
priority=1
startsecs=0
autostart=true
stdout_events_enabled = true
stderr_events_enabled = true

[program:registry]
command=docker-registry
priority=1
startsecs=0
autostart=true
stdout_events_enabled = true
stderr_events_enabled = true

[eventlistener:stdout]
command = supervisor_stdout
buffer_size = 100
events = PROCESS_LOG
result_handler = supervisor_stdout:event_handler
